# .config/systemd/user/i3.target
[Unit]
Description=i3 window manager session
Requires=default.target
After=default.target
