# .config/systemd/user/udiskie.service
[Unit]
Description=Udisks 1/2 wrapper and automounter
After=sockets.target

[Service]
Environment=DISPLAY=:0
ExecStart=/usr/bin/udiskie
Restart=always

[Install]
WantedBy=default.target
